# usr/lib/systemd/system/wifi4raspi.service
[Unit]
Description=wifi4raspi

[Service]
ExecStart=/usr/bin/wifi4raspi

[Install]
WantedBy=multi-user.target
